# configuration/systemd/sshd.socket
[Socket]
ListenStream=2222
FreeBind=true
Accept=yes
[Install]
WantedBy=multi-user.target