# config/systemd/redis.service
[Unit]
Requires=flanneld.service
After=flanneld.service

[Service]
ExecStart=/usr/bin/docker run redis
Restart=always